# instance/scripts/bacalhau.service
[Unit]
Description=Bacalhau Node Service
After=network-online.target docker.service setup-config.service
Requires=docker.service setup-config.service
Wants=network-online.target

[Service]
Type=forking
User=ubuntu
WorkingDirectory=/opt/uploaded_files/scripts
Restart=always
RestartSec=10
StandardOutput=append:/opt/startup.log
StandardError=append:/opt/startup.log

# Wait for Docker to be ready
ExecStartPre=/bin/bash -c 'until docker version > /dev/null 2>&1; do echo "Waiting for Docker..."; sleep 2; done'

# Start Bacalhau using docker compose in detached mode
ExecStart=/usr/bin/docker compose -f /opt/uploaded_files/scripts/docker-compose-bacalhau.yaml up -d
ExecStop=/usr/bin/docker compose -f /opt/uploaded_files/scripts/docker-compose-bacalhau.yaml down

[Install]
WantedBy=multi-user.target default.target