# instance/scripts/sensor.service
[Unit]
Description=Sensor Data Generator Service
After=network-online.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/sensor
Restart=always
RestartSec=10
StandardOutput=append:/opt/sensor/logs/sensor.log
StandardError=append:/opt/sensor/logs/sensor.log

# Wait for Docker to be ready
ExecStartPre=/bin/bash -c 'until docker version > /dev/null 2>&1; do echo "Waiting for Docker..."; sleep 2; done'

# Start sensor using docker-compose
ExecStart=/usr/bin/docker compose -f /opt/sensor/docker-compose.yml up
ExecStop=/usr/bin/docker compose -f /opt/sensor/docker-compose.yml down

[Install]
WantedBy=multi-user.target default.target
