# instance/scripts/setup-config.service
[Unit]
Description=Setup Configuration for Bacalhau and Sensor Services
After=network-online.target cloud-init.target bacalhau-startup.service
Requires=bacalhau-startup.service
Before=bacalhau.service sensor-generator.service
Wants=network-online.target bacalhau.service sensor-generator.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu
RemainAfterExit=yes
StandardOutput=append:/opt/startup.log
StandardError=append:/opt/startup.log

# Create necessary directories
ExecStartPre=/bin/bash -c 'sudo mkdir -p /bacalhau_node /bacalhau_data /opt/sensor/{config,data,logs,exports}'
ExecStartPre=/bin/bash -c 'sudo chown -R ubuntu:ubuntu /bacalhau_node /bacalhau_data /opt/sensor'

# Copy configuration files
ExecStart=/bin/bash -c 'if [ -f /opt/uploaded_files/config/bacalhau-config.yaml ]; then sudo cp /opt/uploaded_files/config/bacalhau-config.yaml /bacalhau_node/config.yaml && sudo chown ubuntu:ubuntu /bacalhau_node/config.yaml; fi'
ExecStart=/bin/bash -c 'if [ -f /opt/uploaded_files/config/sensor-config.yaml ]; then cp /opt/uploaded_files/config/sensor-config.yaml /opt/sensor/config/; fi'

# Generate node identity if it doesn't exist
ExecStart=/bin/bash -c 'if [ ! -f /opt/sensor/config/node_identity.json ] && [ -f /opt/uploaded_files/scripts/generate_node_identity.py ]; then /usr/bin/python3 /opt/uploaded_files/scripts/generate_node_identity.py; fi'

# Create node-info directory for bacalhau
ExecStart=/bin/bash -c 'mkdir -p /bacalhau_node/node-info && echo "Node configured at $(date)" > /bacalhau_node/node-info/setup.txt'

[Install]
WantedBy=multi-user.target default.target