# examples/python-api/.spot/services/api.service
[Unit]
Description=FastAPI Application
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/api
Environment="PATH=/opt/api/venv/bin"
ExecStart=/opt/api/venv/bin/gunicorn main:app -c /opt/api/gunicorn.conf.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fastapi

[Install]
WantedBy=multi-user.target